# base/android/supervisord.conf
[supervisord]
nodaemon=false

[program:xvfb]
priority=10
command=Xvfb :0 -screen 0 1024x768x16

[program:blackbox]
priority=20
command=blackbox
environment=DISPLAY=":0",HOME="/home/user"

[program:x11vnc]
priority=30
command=x11vnc -display :0 -xkb

[program:novnc]
priority=40
command=/opt/noVNC/utils/launch.sh

[program:android_emulator]
priority=0
command=emulator64-arm -avd myandroid
environment=DISPLAY=":0",HOME="/home/user"